# scripts/resources/stellar-history.conf
start on runlevel [2345]
stop on runlevel [06]

setuid stellar
setgid stellar

respawn

chdir /opt/stellar/stellar-core

exec /usr/bin/python -m SimpleHTTPServer
