# scripts/resources/stellar-core.conf
start on runlevel [2345]
stop on runlevel [06]

setuid stellar
setgid stellar

limit nofile 8192 10240

respawn

chdir /opt/stellar/home

exec /opt/stellar/stellar-core/bin/start
